# relay/systemd/bartleby-relay.service
[Unit]
Description=Bartleby hosted API relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bartleby
Group=bartleby
WorkingDirectory=/opt/notique/bartleby-relay
Environment=NODE_ENV=production
EnvironmentFile=/etc/bartleby-relay.env
ExecStart=/usr/bin/node /opt/notique/bartleby-relay/src/server.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/notique/bartleby-relay

[Install]
WantedBy=multi-user.target
